# logitech-wireless-device@.target
[Unit]
Description=Associated services for Logitech wireless device "%i"
