# nm-vpn-id@.target
[Unit]
Description=Associated services for VPN connection "%i"
