# nm-wifi-ssid@.target
[Unit]
Description=Associated services for WiFi network "%i"
